# backup_checkpoint_20251127_084901/playbooks/templates/prometheus-podman.service.j2
# Systemd service for Prometheus Podman container
# Generated by Ansible - DO NOT EDIT MANUALLY

[Unit]
Description=Prometheus Monitoring Container
Documentation=https://prometheus.io/docs/introduction/overview/
After=network-online.target podman.service
Wants=network-online.target
Requires=podman.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/podman start {{ prometheus_container_name }}
ExecStop=/usr/bin/podman stop -t 10 {{ prometheus_container_name }}
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target

