# debian/prerm
#!/bin/bash
# Pre-removal script for InfraMonitor

set -e

APP_NAME="inframonitor"
SERVICE_NAME="${APP_NAME}.service"

# Stop and disable service before removal
if systemctl is-active --quiet "${SERVICE_NAME}" 2>/dev/null; then
    systemctl stop "${SERVICE_NAME}" >/dev/null 2>&1 || true
fi

if systemctl is-enabled --quiet "${SERVICE_NAME}" 2>/dev/null; then
    systemctl disable "${SERVICE_NAME}" >/dev/null 2>&1 || true
fi

exit 0

